# config/webapp.service
[Unit]
Description=Web Application

[Service]
ExecStart=/opt/webapp/bin/webapp
User=csye6225
Group=csye6225

[Install]
WantedBy=multi-user.target
